# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/run/supervisord.pid
childlogdir=/var/log/supervisor

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[unix_http_server]
file=/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ── nginx (reverse proxy + frontend static files) ─────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/tenaos/nginx.log
stderr_logfile=/var/log/tenaos/nginx.err.log

; ── MariaDB (OpenMRS database) ───────────────────────────────────────────
[program:mariadb]
command=/usr/sbin/mariadbd --user=mysql --datadir=/var/lib/mysql --socket=/run/mysqld/mysqld.sock
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/tenaos/mariadb.log
stderr_logfile=/var/log/tenaos/mariadb.err.log

; ── OpenMRS DB bootstrap (one-shot, creates the `openmrs` database + user) ─
[program:init-databases]
command=/opt/tenaos/openmrs-scripts/init-databases.sh
autostart=true
autorestart=false
startsecs=0
priority=25
exitcodes=0
stdout_logfile=/var/log/tenaos/init-databases.log
stderr_logfile=/var/log/tenaos/init-databases.err.log

; ── OpenMRS / Tomcat (long-running) ──────────────────────────────────────
[program:openmrs]
command=/opt/tenaos/openmrs-scripts/run-openmrs.sh
user=tenaos
autostart=true
autorestart=true
priority=30
startsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/tenaos/openmrs.log
stderr_logfile=/var/log/tenaos/openmrs.err.log

; ── Qdrant vector store ──────────────────────────────────────────────────
[program:qdrant]
command=/usr/local/bin/qdrant --config-path /qdrant/config/config.yaml
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/tenaos/qdrant.log
stderr_logfile=/var/log/tenaos/qdrant.err.log

; ── One-shot snapshot restore (no-op if collections already populated) ───
[program:qdrant-restore]
command=/usr/local/bin/tenaos-restore-qdrant
autostart=true
autorestart=false
startsecs=0
priority=45
exitcodes=0
stdout_logfile=/var/log/tenaos/qdrant-restore.log
stderr_logfile=/var/log/tenaos/qdrant-restore.err.log

; ── llama.cpp serving Gemma 4 E4B BF16 GGUF ──────────────────────────────
[program:llama-server]
command=/usr/local/bin/tenaos-start-llama
autostart=true
autorestart=true
priority=50
stdout_logfile=/var/log/tenaos/llama-server.log
stderr_logfile=/var/log/tenaos/llama-server.err.log

; ── KB-Guidelines daemon (loads EmbedGemma in-process) ───────────────────
[program:kb-guidelines]
command=/usr/local/bin/tenaos-start-kb guidelines 4276 who_msf_guidelines
autostart=true
autorestart=true
priority=60
stdout_logfile=/var/log/tenaos/kb-guidelines.log
stderr_logfile=/var/log/tenaos/kb-guidelines.err.log

; ── KB-CIEL daemon (loads EmbedGemma in-process) ─────────────────────────
[program:kb-ciel]
command=/usr/local/bin/tenaos-start-kb ciel 4277 ciel_concepts
autostart=true
autorestart=true
priority=60
stdout_logfile=/var/log/tenaos/kb-ciel.log
stderr_logfile=/var/log/tenaos/kb-ciel.err.log

; ── TenaAgent (Python orchestrator) ──────────────────────────────────────
[program:tena-agent]
command=/usr/local/bin/tenaos-start-tena-agent
user=tenaos
autostart=true
autorestart=true
priority=70
stdout_logfile=/var/log/tenaos/tena-agent.log
stderr_logfile=/var/log/tenaos/tena-agent.err.log
