# TenaOS-Backend/configs/supervisor/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:mariadb]
command=/usr/sbin/mariadbd --datadir=/var/lib/mysql --socket=/run/mysqld/mysqld.sock --pid-file=/run/mysqld/mysqld.pid
user=mysql
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/mariadb.log
stderr_logfile=/var/log/supervisor/mariadb.log
startsecs=5

[program:init-databases]
command=/opt/tenaos/scripts/init-databases.sh
autostart=true
autorestart=false
startsecs=0
priority=20
exitcodes=0
stdout_logfile=/var/log/supervisor/init-databases.log
stderr_logfile=/var/log/supervisor/init-databases.log

[program:openmrs]
command=/opt/tenaos/scripts/run-openmrs.sh
user=tenaos
autostart=true
autorestart=true
priority=30
startsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/openmrs.log
stderr_logfile=/var/log/supervisor/openmrs.log
